# deploy/com.tradingbot.paper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  com.tradingbot.paper — Gölge paper scheduler (F5-B1) günlük EOD işi.

  NE YAPAR: BIST kapanışı (18:00 Istanbul) + grace sonrası her İŞ GÜNÜ bir kez
  koşar: yfinance EOD veri güncelle (--refresh) + günlük döngü (--cycle). Tatil/
  hafta sonu günleri scheduler kendi içinde bist_calendar ile atlar (zararsız koşar).

  KURULUM (operatör, README/OPERATOR_GUIDE §Başlat):
    1) Bir kez bootstrap:  .venv/bin/python main.py --config config/regime_core.yaml --bootstrap
    2) cp deploy/com.tradingbot.paper.plist ~/Library/LaunchAgents/
    3) launchctl load ~/Library/LaunchAgents/com.tradingbot.paper.plist

  ZAMANLAMA GEREKÇESİ (K5 grace uyumu): BIST sürekli kapanış 18:00 Istanbul; K5
  bar-finallik grace = 3600s (bar_close_grace_sec) → bar 19:00'da FİNAL olur.
  Sinyal yalnız FİNAL bardan hesaplanır; erken koşarsa cycle provisional görür ve
  observe'da işlem yapmaz / active'de yürütmeyi son final güne erteler (zarar yok,
  ama o gün atlanır). Bu yüzden iş 19:30'da koşar → 19:05 alt-sınırına 25 dk marj.
  Türkiye kalıcı UTC+3 (DST YOK), dolayısıyla 19:30 yerel yıl boyu sabit güvenli.

  NOT: Yol/kullanıcı-adı Mac Mini'ye göre; gerekiyorsa düzenle. StartCalendarInterval
  SİSTEM yerel saatini kullanır — Mac Istanbul (UTC+3) saatinde varsayılıyor; farklı
  tz'de Hour'u >=19:05 Istanbul karşılığına ayarla (kod tarafı to_istanbul ile tz-güvenli,
  yalnız TETİKLEME saati sisteme bağlı). 'real'e geçiş YOK — yalnızca paper (Durma Noktası 2).
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tradingbot.paper</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/sano/trading-bot/.venv/bin/python</string>
        <string>/Users/sano/trading-bot/main.py</string>
        <string>--config</string>
        <string>/Users/sano/trading-bot/config/regime_core.yaml</string>
        <string>--refresh</string>
        <string>--cycle</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/sano/trading-bot</string>

    <!-- Her iş günü 19:30 yerel = kapanış 18:00 + 1h grace (19:00 final) + 30dk marj.
         Cumartesi/Pazar hariç; BIST tatilleri scheduler'da bist_calendar ile atlanır. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>19</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>19</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>19</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>19</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>19</integer><key>Minute</key><integer>30</integer></dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/sano/trading-bot/runtime/paper/logs/paper.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/sano/trading-bot/runtime/paper/logs/paper.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
